# scripts/packaging/systemd/cairo-web.service
[Unit]
Description=Cairo web agent

[Service]
ExecStart=/usr/local/bin/cairo-web
Restart=on-failure

[Install]
WantedBy=default.target
